# src/lab3/capture.service
[Unit]
Description=fuck

[Service]
ExecStart=/root/lab3/start_capture.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
